# clipboard-daemon.service
[Unit]
Description=Clipboard History Daemon
Documentation=https://github.com/yourusername/clipboard-history
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/clipboard-daemon
Restart=on-failure
RestartSec=5

# Security hardening
PrivateTmp=true
NoNewPrivileges=true

# Environment
Environment="DISPLAY=:0"
Environment="XAUTHORITY=%h/.Xauthority"

[Install]
WantedBy=default.target